# supervisord.conf
[supervisord]
nodaemon=true

[program:apacheds]
command=/opt/apacheds/bin/startApacheDS
